# testingScripts/crontab
# NOTE: Embedding the command sent to derecho via ssh in double-quotes
# ensures proper remote argument parsing.

# Run quick tests every night.
34 00 * * * ssh derecho "/glade/work/ewinter/mage_testing/derecho/kaiju/ewinter-derecho_testing_fixes/kaiju/testingScripts/run_mage_tests.sh -v 'buildTest.py -lv,ICtest.py -lv'" >> /glade/work/ewinter/mage_testing/derecho/logs/nightly-tests.out 2>&1

# Run tests requiring PBS jobs on Monday nights.
34 00 * * 1 ssh derecho "/glade/work/ewinter/mage_testing/derecho/kaiju/ewinter-derecho_testing_fixes/kaiju/testingScripts/run_mage_tests.sh -v 'pyunitTest.py -lv,intelChecks.py -lv,unitTest.py -lv'" >> /glade/work/ewinter/mage_testing/derecho/logs/weekly-tests.out 2>&1
